# cohere-stt/cohere-stt-api.service
[Unit]
Description=Cohere STT API Server (local speech-to-text)
After=network.target

[Service]
Type=simple
User=zhugehyuk
WorkingDirectory=/home/zhugehyuk/2lab.ai/skills/cohere-stt/scripts
ExecStart=/home/zhugehyuk/.cohere-stt-venv/bin/python transcribe_api.py --port 8787 --preload
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

# GPU access
SupplementaryGroups=video render

[Install]
WantedBy=multi-user.target
